# templates/linux/stud.init.conf
#!upstart
description "starting stud"
author      "comet"

start on runlevel [2345]
stop on runlevel [06]

respawn
limit nofile 65536 65536

script
  stud --config=/opt/stud/stud.conf
end script